# server/deploy/gobang-server.service
[Unit]
Description=Gobang Game Server
After=network.target

[Service]
Type=simple
User=gobang
Group=gobang
ExecStart=/usr/local/bin/gobang-server
WorkingDirectory=/var/lib/gobang-server
Restart=on-failure
RestartSec=5

Environment=GOBANG_DATA_DIR=/var/lib/gobang-server
Environment=GOBANG_PORT=3001
Environment=GOBANG_HOST=0.0.0.0

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/gobang-server
PrivateTmp=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=gobang-server

[Install]
WantedBy=multi-user.target
